# DEBIAN/postinst
#!/bin/sh
set -e
#DEBHELPER#

case "$1" in
  configure)
    # Allow any user to run chvt 2 (for the desktop shortcut)
    # Allow the moonpi service user to write to tty1 (for the exit-to-console message)
    printf 'ALL ALL=(root) NOPASSWD: /usr/bin/chvt 2\nmoonpi ALL=(root) NOPASSWD: /usr/bin/tee /dev/tty1\n' > /etc/sudoers.d/moonpi
    chmod 440 /etc/sudoers.d/moonpi

    # deb-systemd-helper does not support template unit instances —
    # systemctl is used directly (see debian/moonpi.lintian-overrides)
    if [ -d /run/systemd/system ]; then
      systemctl daemon-reload                        || true
      systemctl enable moonpi@tty2.service           || true
      systemctl enable switch-to-tty2.service        || true
      systemctl start  moonpi@tty2.service           || true
      systemctl start  switch-to-tty2.service        || true
    else
      deb-systemd-helper enable moonpi@tty2.service    || true
      deb-systemd-helper enable switch-to-tty2.service || true
    fi
    ;;
esac

exit 0